# supervisor/syncthing.conf
[program:syncthing]
command = /app/syncthing/syncthing -no-browser -home="/config/syncthing"
directory = /data
autorestart = True
user = abc
environment = STNORESTART="1", HOME="/data"