# supervisor/syncthing.conf
[program:syncthing]
command=syncthing -home=/config/syncthing -no-browser -no-restart
user=abc